# noip2.service
[Unit]
Description=NOIp Dynamic ip

[Service]
Type=forking
ExecStart=/usr/local/bin/noip2

[Install]
WantedBy=multi-user.target
